# cron/flightfinder.cron
# flightFinderWeb — scrape Google Flights for direct flights every hour
# Install with: (crontab -l 2>/dev/null; cat cron/flightfinder.cron) | crontab -
# Remove with:  crontab -e  (and delete the line)

0 * * * * /Users/eeisenberg/work/src/flightFinderWeb/cron/run.sh >> /var/log/flightfinder.log 2>&1
